# sample/crontab
*/1 * * * * /usr/bin/sh ${HOME}/.local/opt/Share_Photo/bin/done.sh
